# root/etc/s6-overlay/s6-rc.d/code-server/run
#!/command/execlineb -P

# s6-notifyoncheck -d -n 300 -w 1000 -c "nc -z 127.0.0.1 8443" 
runuser vscode -- -lc "exec \"$@\"" -
code-server
# --bind-addr 0.0.0.0:8443
# --user-data-dir /config/data
# --extensions-dir /config/extensions
--disable-telemetry
/mnt
