# root/etc/s6-overlay/s6-rc.d/jupyterlab/run
#!/command/execlineb -P

runuser vscode -l -- -c "exec \"$@\"" -
jupyter lab
--no-browser
/mnt
