# root/etc/s6-overlay/s6-rc.d/init-script/run
#!/command/with-contenv bash

if [[ ! -x "$INIT_FILE" || -f /etc/init_finished ]] ; then
    exit 0
fi

touch /etc/init_finished

runuser vscode -l -- $INIT_FILE
